# mrl-memory.service
[Unit]
Description=MRL Memory System API
After=network.target

[Service]
Type=simple
User=your_user
WorkingDirectory=/path/to/manaos_integrations
EnvironmentFile=/path/to/manaos_integrations/.env
ExecStart=/usr/bin/python3 /path/to/manaos_integrations/mrl_memory_integration.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
